# docker-aws-cli/crontabs/aws
# min	hour	day	month	weekday	command
0	*	*	*	*	/scripts/mysqldump.sh >> /var/log/mysqldump.log
0	5	*	*	*	/scripts/backupfileuserdata.sh >> /var/log/backupfileuserdata.log
