# docker-mikrotik/supervisord.conf
[supervisord]
nodaemon=true

[program:qemu]
command=/entrypoint.sh
autorestart=true

[program:novnc]
command=/opt/noVNC/utils/launch.sh --vnc localhost:5900 --listen 8080
autorestart=true
